# deploy/tripleader.service
[Unit]
Description=MUUC Trip Organizer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tripleader
EnvironmentFile=/opt/tripleader/.env
ExecStart=/opt/tripleader/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8999
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
